# firmware/85-n2pkvna.rules.in
SUBSYSTEM=="usb", ATTR{idVendor}=="0547", ATTR{idProduct}=="1004", ACTION=="add", RUN+="/sbin/fxload -t fx2 -I @prefix@/lib/firmware/n2pkvna/n2pkvna-v4.hex -D %N"
SUBSYSTEMS=="usb", ATTRS{idVendor}=="0547", ATTRS{idProduct}=="1005", ACTION=="add", GROUP="n2pkvna", MODE="0660", SYMLINK+="n2pkvna%n"
SUBSYSTEM=="usb", ATTR{idVendor}=="0547", ATTR{idProduct}=="1008", ACTION=="add", RUN+="/sbin/fxload -t fx2 -I @prefix@/lib/firmware/n2pkvna/n2pkvna-v4-500.hex -D %N"
SUBSYSTEMS=="usb", ATTRS{idVendor}=="0547", ATTRS{idProduct}=="1009", ACTION=="add", GROUP="n2pkvna", MODE="0660", SYMLINK+="n2pkvna%n"
SUBSYSTEM=="usb", ATTR{idVendor}=="0547", ATTR{idProduct}=="100a", ACTION=="add", RUN+="/sbin/fxload -t fx2 -I @prefix@/lib/firmware/n2pkvna/n2pkvna-v5.hex -D %N"
SUBSYSTEMS=="usb", ATTRS{idVendor}=="0547", ATTRS{idProduct}=="100b", ACTION=="add", GROUP="n2pkvna", MODE="0660", SYMLINK+="n2pkvna%n"
SUBSYSTEM=="usb", ATTR{idVendor}=="0547", ATTR{idProduct}=="100c", ACTION=="add", RUN+="/sbin/fxload -t fx2 -I @prefix@/lib/firmware/n2pkvna/n2pkvna-v5-500.hex -D %N"
SUBSYSTEMS=="usb" ATTRS{idVendor}=="0547", ATTRS{idProduct}=="100d", ACTION=="add", GROUP="n2pkvna", MODE="0660", SYMLINK+="n2pkvna%n"
